# templates/default/upstart/centos-component.conf.erb
#!upstart
description "Prometheus platform: <%= @comp %>"

start on started network
stop on stopping network

respawn
respawn limit unlimited

chdir <%= @path %>
exec su -s /bin/sh -c 'exec "$0" "$@" >>/var/log/prometheus/<%= @comp %>.log 2>&1' <%= @user %> -- <%= @cli %>
